# conf/google-cloud-storage-sync.cron
43 * * * * root gsutil cp /opt/betscrape/logs/* gs://betscrape-api-logs/
