# conf/google-cloud-storage-sync.cron
10 * * * * * gsutil cp /opt/betscrape/logs/* gs://betscrape-api-logs/
